# data/etc/systemd/system/immutable-update.service
[Unit]
Description=Immutable root update (stage snapshot via nspawn)
Documentation=file:///etc/immutable-update.conf
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/immutable-update
TimeoutStartSec=3600
# Keep the nspawn container's child processes confined
KillMode=mixed
# Journald will capture stdout/stderr alongside the log file
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
