# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        systemctl daemon-reload 2>/dev/null || true
        systemctl enable immutable-update.timer 2>/dev/null || true
        systemctl enable machine-id-persist.service 2>/dev/null || true

        # On upgrades (not first install), regenerate initramfs and GRUB
        if [ -n "$2" ]; then
            dracut --regenerate-all --force 2>/dev/null || true
            update-grub 2>/dev/null || true
        fi
        ;;
esac

#DEBHELPER#
